# service/rfmqtt.service
[Unit]
Description=RF 433MHz to MQTT bridge

[Service]
ExecStart=/usr/local/bin/rfmqtt
Restart=on-failure

[Install]
WantedBy=default.target

